# templates/launchd/hermes-dashboard.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>hermes-dashboard</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>exec "$HOME/.local/bin/hermes" dashboard --host 127.0.0.1 --port {{HERMES_DASHBOARD_PORT}} --no-open</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/hermes-dashboard.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/hermes-dashboard.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HERMES_DASHBOARD_PORT</key>
    <string>{{HERMES_DASHBOARD_PORT}}</string>
  </dict>
</dict>
</plist>
